# config/supervisord.conf
[inet_http_server] 
port=127.0.0.1:9001 

[rpcinterface:supervisor] 
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface 

[supervisorctl] 
serverurl=http://127.0.0.1:9001 

[supervisord] 

[program:doc_chat] 
environment = PYTHONUNBUFFERED=1 
user=root 
directory=/code 
command=gunicorn -w 2 -k uvicorn.workers.UvicornWorker app:app 
autostart=true 
autorestart=true 
stderr_logfile=/var/log/doc_chat/doc_chat.err.log 
stdout_logfile=/var/log/doc_chat/doc_chat.out.log